# deploy/proofofwork-postgres-logical-backup.service
[Unit]
Description=ProofOfWork PostgreSQL logical backup
Requires=postgresql@16-main.service
After=postgresql@16-main.service
RequiresMountsFor=/data/proofofwork-postgres-backups

[Service]
Type=oneshot
User=postgres
Group=postgres
UMask=0077
ExecStart=/usr/local/sbin/proofofwork-postgres-logical-backup
TimeoutStartSec=2h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/proofofwork-postgres-backups/logical
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
