# deploy/proofofwork-node-release-prune.timer
[Unit]
Description=Daily managed ProofOfWork node release retention

[Timer]
OnCalendar=daily
RandomizedDelaySec=15min
FixedRandomDelay=true
Persistent=true
Unit=proofofwork-node-release-prune.service

[Install]
WantedBy=timers.target
